# deploy/sunyudx-api.service
[Unit]
Description=sunyuDX-flow API Server
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=exec
User=sunyudx
Group=sunyudx
WorkingDirectory=/opt/sunyudx-flow/backend

# Environment
Environment="PATH=/opt/sunyudx-flow/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/sunyudx-flow/backend/.env.production

# Start command (Uvicorn with multiple workers)
ExecStart=/opt/sunyudx-flow/backend/venv/bin/uvicorn main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --loop uvloop \
    --http httptools \
    --access-log \
    --log-level info

# Restart policy
Restart=always
RestartSec=5

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/sunyudx-flow/uploads /var/log/sunyudx

# Logging
StandardOutput=append:/var/log/sunyudx/api.log
StandardError=append:/var/log/sunyudx/api-error.log

[Install]
WantedBy=multi-user.target
